# serverstat.service
[Unit]
Description=Server-Stat Monitoring Service
After=network.target

[Service]
Type=simple
ExecStart=/opt/serverstat/serverstat -config /opt/serverstat/config.yaml
Restart=always
RestartSec=5
User=root
WorkingDirectory=/opt/serverstat

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/serverstat
PrivateTmp=true

[Install]
WantedBy=multi-user.target
